# deploy/hetzner/dressai.service
[Unit]
Description=DressAI API (dressai-api/server.js)
After=network.target

[Service]
Type=simple
User=dressai
Group=dressai
WorkingDirectory=/home/dressai/getdressai.com/dressai-api
Environment=NODE_ENV=production
EnvironmentFile=/home/dressai/getdressai.com/dressai-api/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
